# snapshots/v7.2_2026-04-04/project/infra/systemd/trading-daemon.service
[Unit]
Description=Trading Execution Daemon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/trading
EnvironmentFile=/etc/trading/trading.env
ExecStart=/opt/trading/.venv/bin/python -m src.operational_daemon --env-file /etc/trading/trading.env
Restart=always
RestartSec=10
StandardOutput=append:/var/log/trading/daemon.log
StandardError=append:/var/log/trading/daemon-error.log

[Install]
WantedBy=multi-user.target
